# zfsdash.service
[Unit]
Description=ZFSdash — ZFS Management Dashboard
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/zfsdash
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zfsdash

[Install]
WantedBy=multi-user.target
